# scripts/crontab
0 0 13 * * sh /usr/scripts/backup.sh
0 0 1 1 * sh /usr/scripts/backup.sh